# supervisor-vncd.conf
[program:vncd]
command=bash -c 'USER=%(ENV__IMG_USERNAME)s vncserver -geometry 1280x800 :0'
autorestart=false
directory=/home/dave
